# easytier-installer/systemd/easytier.service
[Unit]
Description=EasyTier - Simple Mesh VPN Service
Documentation=https://github.com/EasyTier/EasyTier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/easytier-core --ipv4 10.145.0.1 --network-name my-network --listeners udp://0.0.0.0:11011
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security settings (مینیمال)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/easytier /var/log

# Environment
Environment=RUST_LOG=info
EnvironmentFile=-/etc/easytier/environment

# Working directory
WorkingDirectory=/etc/easytier

# PID file
PIDFile=/run/easytier.pid

[Install]
WantedBy=multi-user.target 